# audio-normalizer.openrc
#!/sbin/openrc-run

name="audio-normalizer"
description="Audio Normalizer service"
command="/opt/audio-normalizer/audio-normalizer.sh"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"
command_user="root"  # Change this to the appropriate user if needed

depend() {
    need localmount
    use net
}

start_pre() {
    checkpath --directory --owner root:root --mode 0755 /run
}